# common/root/etc/s6-overlay/s6-rc.d/init-desktop-files/up
with-contenv /etc/s6-overlay/s6-rc.d/init-desktop-files/start.sh
