# common/root/etc/s6-overlay/s6-rc.d/svc-desktop/run
#!/command/execlineb -S0
with-contenv /etc/s6-overlay/s6-rc.d/svc-desktop/start.sh
