# common/root/etc/s6-overlay/s6-rc.d/svc-custom/run
#!/command/execlineb -S0
with-contenv /etc/helios/services.d/custom.sh
